# packages/atlas-system-daemon/lib/systemd/system/atlas-system-daemon.service
[Unit]
Description=Atlas System Daemon
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/atlas/system_daemon.py
RuntimeDirectory=atlas
Restart=on-failure
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/run/atlas /srv/atlas/logs
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=multi-user.target
