# config/includes.chroot/etc/skel/.config/autostart/atlas-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Atlas Wallpaper
Comment=Set the Atlas glowing-triangle desktop wallpaper
Exec=/usr/lib/atlas/set-wallpaper.sh
OnlyShowIn=XFCE;
X-XFCE-Autostart-Override=true
X-GNOME-Autostart-enabled=true
Hidden=false
